# bitcoin.service
[Unit]
Description=Bitcoin Core Daemon
After=network.target

[Service]
Type=forking
User=bitcoin
Group=bitcoin

# Main process
ExecStart=/usr/local/bin/bitcoind -conf=/etc/bitcoin/bitcoin.conf -datadir=/var/lib/bitcoin -daemon
ExecStop=/usr/local/bin/bitcoin-cli -conf=/etc/bitcoin/bitcoin.conf -datadir=/var/lib/bitcoin stop

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bitcoind

[Install]
WantedBy=multi-user.target
